# templates/systemd/system/mediaplayer.service
[Unit]
Description=Mediaplayer
Requires=mpd.service
After=mpd.service

[Service]
Type=simple
User={{ mediaplayer_user_name }}
WorkingDirectory={{ mediaplayer_install_location }}/mediaplayer
ExecStart=/usr/bin/java -jar mediaplayer.jar
Restart=on-abort

[Install]
WantedBy=multi-user.target
